# cron/cron.d

*/2 * * * * /usr/local/bin/python /app/backup.py >> /dev/null
