# cron/cron.d

*/2 * * * * python /app/backup.py > /dev/null
